# projet_08-supervisor.conf
[program:projet_08-gunicorn]
command = /usr/local/bin/gunicorn Projet_08.wsgi:application
user = manu
directory = /home/manu/Projet_08
autostart = true
autorestart = true
